# system/chrome.service
[Unit]
Description=Chrome-Browser

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/kiosk.sh

[Install]
WantedBy=multi-user.target
